# auxiliary/wiki_monkey.service
[Unit]
Description=Wiki Monkey's database server

[Service]
ExecStart=/usr/bin/wiki-monkey

[Install]
WantedBy=default.target
